# systemd/.config/systemd/user/unclutter.service
[Unit]
Description=Erase mouse cursor

[Service]
ExecStart=/usr/bin/unclutter -idle 1.0 -root

[Install]
WantedBy=graphical.target
